# docker/supervisord.conf
# ╔══════════════════════════════════════════════════════════════════════╗
# ║  NovaPanel Supervisord Configuration                                ║
# ║                                                                     ║
# ║  Manages all service daemons inside the Docker container.           ║
# ║  Priority order: databases → caches → web/php → mail/dns/ftp →     ║
# ║  security → panel                                                   ║
# ╚══════════════════════════════════════════════════════════════════════╝

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ─── Priority 10: Databases ─────────────────────────────────────────────

[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql
autostart=true
autorestart=true
priority=10
startsecs=10
startretries=5
stdout_logfile=/var/log/supervisor/mariadb.log
stderr_logfile=/var/log/supervisor/mariadb-error.log

[program:postgresql]
command=bash -c "chown -R postgres:postgres /var/lib/postgresql /var/run/postgresql && exec su - postgres -c '/usr/lib/postgresql/%(ENV_PG_VERSION)s/bin/postgres -D /var/lib/postgresql/%(ENV_PG_VERSION)s/main -c config_file=/etc/postgresql/%(ENV_PG_VERSION)s/main/postgresql.conf'"
environment=PG_VERSION="16"
autostart=true
autorestart=true
priority=10
startsecs=10
startretries=5
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql-error.log

# ─── Priority 15: Caching ───────────────────────────────────────────────

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
priority=15
startsecs=5
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-error.log

# ─── Priority 20: Web Server & PHP-FPM ──────────────────────────────────

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log

[program:php8.1-fpm]
command=/usr/sbin/php-fpm8.1 --nodaemonize --fpm-config /etc/php/8.1/fpm/php-fpm.conf
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/php8.1-fpm.log
stderr_logfile=/var/log/supervisor/php8.1-fpm-error.log

[program:php8.2-fpm]
command=/usr/sbin/php-fpm8.2 --nodaemonize --fpm-config /etc/php/8.2/fpm/php-fpm.conf
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/php8.2-fpm.log
stderr_logfile=/var/log/supervisor/php8.2-fpm-error.log

[program:php8.3-fpm]
command=/usr/sbin/php-fpm8.3 --nodaemonize --fpm-config /etc/php/8.3/fpm/php-fpm.conf
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/var/log/supervisor/php8.3-fpm.log
stderr_logfile=/var/log/supervisor/php8.3-fpm-error.log

# ─── Priority 25: Mail Stack ────────────────────────────────────────────

[program:postfix]
command=bash -c "postfix start-fg || exec /usr/lib/postfix/sbin/master -d"
autostart=true
autorestart=true
priority=25
startsecs=5
stdout_logfile=/var/log/supervisor/postfix.log
stderr_logfile=/var/log/supervisor/postfix-error.log

[program:dovecot]
command=/usr/sbin/dovecot -F
autostart=true
autorestart=true
priority=25
startsecs=5
stdout_logfile=/var/log/supervisor/dovecot.log
stderr_logfile=/var/log/supervisor/dovecot-error.log

# ─── Priority 25: DNS ───────────────────────────────────────────────────

[program:bind9]
command=/usr/sbin/named -f -u bind -g
autostart=true
autorestart=true
priority=25
startsecs=5
stdout_logfile=/var/log/supervisor/bind9.log
stderr_logfile=/var/log/supervisor/bind9-error.log

# ─── Priority 25: FTP ───────────────────────────────────────────────────

[program:proftpd]
command=/usr/sbin/proftpd -n
autostart=true
autorestart=true
priority=25
startsecs=5
stdout_logfile=/var/log/supervisor/proftpd.log
stderr_logfile=/var/log/supervisor/proftpd-error.log

# ─── Priority 30: Security ──────────────────────────────────────────────

[program:fail2ban]
command=/usr/bin/fail2ban-server -f
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/var/log/supervisor/fail2ban.log
stderr_logfile=/var/log/supervisor/fail2ban-error.log

# ─── Priority 30: Cron ──────────────────────────────────────────────────

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=30
startsecs=3
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron-error.log

# ─── Priority 40: NovaPanel API ─────────────────────────────────────────
# Waits for MariaDB and Redis before starting

[program:novapanel]
command=bash -c "echo '[novapanel] Waiting for dependencies...'; for i in $(seq 1 60); do if mariadb -u root -e 'SELECT 1' >/dev/null 2>&1 && redis-cli ping >/dev/null 2>&1; then echo '[novapanel] Dependencies ready!'; break; fi; echo '[novapanel] Waiting... ('$i'/60)'; sleep 1; done; exec /usr/bin/node /opt/novapanel/apps/api/dist/index.js"
directory=/opt/novapanel/apps/api
user=novapanel
environment=NODE_ENV="production",PORT="8732",HOST="0.0.0.0",DB_PATH="/var/lib/novapanel/novapanel.db",REDIS_URL="redis://127.0.0.1:6379",PANEL_URL="%(ENV_PANEL_URL)s",ADMIN_EMAIL="%(ENV_ADMIN_EMAIL)s",LE_EMAIL="%(ENV_LE_EMAIL)s",MAIL_HOSTNAME="%(ENV_MAIL_HOSTNAME)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",VHOSTS_ROOT="/var/www/vhosts",BACKUP_DIR="/var/lib/novapanel/backups"
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=5
stdout_logfile=/var/log/supervisor/novapanel.log
stderr_logfile=/var/log/supervisor/novapanel-error.log
